# ovs-agent/systemd/lbaas_ovs_agent.service
[Unit]
Description=LBaaS OVS Agent (FastAPI)
After=network.target docker.service
Requires=docker.service

[Service]
WorkingDirectory=/opt/lbaasovsagent
ExecStart=/opt/lbaasovsagent/venv/bin/lbaas-ovs-agent
User=lbaasovsagent
Group=lbaasovsagent
Restart=always
RestartSec=2
Environment="PYTHONUNBUFFERED=1"
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/lbaasovsagent
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target